# root/etc/systemd/system/backrest.service
[Service]
ExecStart=/usr/bin/backrest
Environment=BACKREST_RESTIC_COMMAND=/usr/bin/restic
Restart=on-failure
[Install]
WantedBy=multi-user.target
