# root/etc/systemd/system/kanata.service
[Service]
WorkingDirectory=/home/ripytide
ExecStart=bash -c kanata
Restart=always

[Install]
WantedBy=multi-user.target
