# root/etc/systemd/system/ydotool.service
[Service]
WorkingDirectory=/
ExecStart=/bin/ydotoold
Restart=always

[Install]
WantedBy=multi-user.target
